# etc/cron.d/redislabs
*/5 * * * *     root    /opt/redislabs/bin/ccs_rdb_rotate.sh 5min 24
